# platform/dracut/99wootc-boot/wootc-attach.service
[Unit]
Description=wootc: attach the NTFS-hosted root.disk before sysroot.mount
Documentation=https://github.com/tuna-os/wootc
DefaultDependencies=no
# THE attach must happen from THIS service, at a point where systemd-udevd is up
# and TRACKING devices, so the loop partitions it creates get promoted into
# dev-disk-by-uuid-<root>.device UNITS that sysroot.mount can wait on. The
# belt-and-braces initqueue hook (removed) attached far too early (t≈1.8s) —
# systemd-udevd was not yet tracking devices, so the by-uuid symlink appeared but
# the .device unit never activated and sysroot.mount timed out anyway.
#
# Ordering is the whole game: after udev is up, before anything waits for the
# root device or mounts sysroot.
#
# Do NOT depend on systemd-udev-settle.service — it is DEPRECATED and, ordered
# ahead of us, it TIMED OUT and delayed this service past the 2-minute
# dev-disk-by-uuid device timeout (proven: sysroot.mount failed while the attach
# was still blocked behind a settle that never completed). The script does its
# own targeted `udevadm trigger --settle` on the loop device instead.
After=systemd-udevd.service systemd-udev-trigger.service
Before=initrd-root-device.target sysroot.mount
Wants=initrd-root-device.target
# Run even in the emergency path so a partial failure is still observable.
ConditionKernelCommandLine=loop

[Service]
Type=oneshot
RemainAfterExit=yes
# If a kernel has no ntfs3, ntfs-3g leaves a FUSE daemon in this service's
# cgroup. initrd cleanup stops the oneshot before switch-root; KillMode=control-
# group would kill that daemon and turn the already-mounted loop root into EIO.
# The shell main process has exited by then, so process mode preserves only the
# backing daemon that the real root still needs.
KillMode=process
# The script is re-entrant (returns 0 once /run/wootc-loop-attached exists) and
# self-diagnosing (writes every exit reason to /dev/console and kmsg).
ExecStart=/usr/bin/bash /usr/lib/wootc/wootc-attach-loop.sh
# Never let a hiccup here wedge the boot; the script reports its own failures.
TimeoutStartSec=90

[Install]
WantedBy=initrd-root-device.target
