# payload/migration/wootc-esp-sync.service
[Unit]
Description=Sync Phase-2 kernel to the Windows ESP after OS updates
# Needs /boot (deployment kernels) and udev by-uuid links settled.
After=local-fs.target systemd-udev-settle.service
ConditionPathExists=/etc/wootc/host-esp.conf

[Service]
Type=oneshot
ExecStart=/var/usrlocal/bin/wootc-esp-sync

[Install]
WantedBy=multi-user.target
